# smf/manifests/cn-agent-setup.xml.in
<?xml version="1.0"?>
<!DOCTYPE service_bundle SYSTEM "/usr/share/lib/xml/dtd/service_bundle.dtd.1">
<!--
    This Source Code Form is subject to the terms of the Mozilla Public
    License, v. 2.0. If a copy of the MPL was not distributed with this
    file, You can obtain one at http://mozilla.org/MPL/2.0/.
-->

<!--
    Copyright (c) 2017, Joyent, Inc.
-->

<service_bundle type="manifest" name="cn-agent-setup">
  <service name="smartdc/agent/cn-agent-setup" type="service" version="@@VERSION@@">

    <create_default_instance enabled="@@ENABLED@@"/>
    <single_instance/>

    <dependency name="network" grouping="require_all" restart_on="error" type="service">
      <service_fmri value="svc:/milestone/network:default"/>
    </dependency>

    <exec_method
      type="method"
      name="start"
      exec="@@ROOT@@/smf/method/cn-agent-setup %m"
      timeout_seconds="600">
      <method_context>
        <method_credential user="root" group="staff"/>
        <method_environment>
          <envvar name="PATH" value="/usr/bin:/usr/sbin"/>
        </method_environment>
      </method_context>
    </exec_method>

    <exec_method type="method" name="stop" exec=":true" timeout_seconds="60" />

    <property_group name="startd" type="framework">
      <propval name="duration" type="astring" value="transient" />
      <propval name="ignore_error" type="astring" value="core,signal" />
    </property_group>

    <property_group name="application" type="application">
    </property_group>

    <stability value="Evolving"/>

    <template>
      <common_name>
        <loctext xml:lang="C">SmartDataCentre Compute Node Agent Setup</loctext>
      </common_name>
    </template>

  </service>
</service_bundle>
